# zscorecrunch/crontab
0 12 * * 1-5 /Users/horse/Projects/ZScoreCrunch/run_daily.sh 2>&1 | tee /tmp/cron_daily.log
